# files/etc/service/vamp-gateway-agent/run
#! /bin/bash

handle() { echo "vamp-gateway-agent/run: got signal"; exit; }
trap handle SIGINT

[[ -e /usr/local/vamp/good2go ]] && exit 0

: "${VAMP_GATEWAY_AGENT_LOGO:=TRUE}"

if [ "$VAMP_GATEWAY_AGENT_LOGO" = "TRUE" ] || [ "$VAMP_GATEWAY_AGENT_LOGO" = "1" ]; then
echo "
██╗   ██╗ █████╗ ███╗   ███╗██████╗      ██████╗  █████╗ ████████╗███████╗██╗    ██╗ █████╗ ██╗   ██╗
██║   ██║██╔══██╗████╗ ████║██╔══██╗    ██╔════╝ ██╔══██╗╚══██╔══╝██╔════╝██║    ██║██╔══██╗╚██╗ ██╔╝
██║   ██║███████║██╔████╔██║██████╔╝    ██║  ███╗███████║   ██║   █████╗  ██║ █╗ ██║███████║ ╚████╔╝
╚██╗ ██╔╝██╔══██║██║╚██╔╝██║██╔═══╝     ██║   ██║██╔══██║   ██║   ██╔══╝  ██║███╗██║██╔══██║  ╚██╔╝
 ╚████╔╝ ██║  ██║██║ ╚═╝ ██║██║         ╚██████╔╝██║  ██║   ██║   ███████╗╚███╔███╔╝██║  ██║   ██║
  ╚═══╝  ╚═╝  ╚═╝╚═╝     ╚═╝╚═╝          ╚═════╝ ╚═╝  ╚═╝   ╚═╝   ╚══════╝ ╚══╝╚══╝ ╚═╝  ╚═╝   ╚═╝
"
fi

: "${VAMP_KEY_VALUE_STORE_PATH?not provided.}"
: "${VAMP_KEY_VALUE_STORE_TYPE?not provided.}"
: "${VAMP_KEY_VALUE_STORE_CONNECTION?not provided.}"
: "${VAMP_PERSISTENCE_STORE_CONNECTION?not provided.}"


for try in {1..9} ; do
  echo "vamp-gateway-agent/run: waiting for syslog to come online..."
  [[ -e /var/run/rsyslogd.pid ]] && break || sleep 1
done


printf "VAMP GATEWAY AGENT: " && cat /usr/local/vamp/version
/usr/bin/confd -version

echo "VAMP_KEY_VALUE_STORE_TYPE       : ${VAMP_KEY_VALUE_STORE_TYPE}"
echo "VAMP_KEY_VALUE_STORE_CONNECTION : ${VAMP_KEY_VALUE_STORE_CONNECTION}"
echo "VAMP_KEY_VALUE_STORE_PATH       : ${VAMP_KEY_VALUE_STORE_PATH}"

mkdir -p /usr/local/vamp/confd/conf.d
mkdir -p /usr/local/vamp/confd/templates

echo "vamp-gateway-agent/run: creating confd configuration and template"
cat <<EOT > /usr/local/vamp/confd/conf.d/workflow.toml
[template]
src = "haproxy.tmpl"
dest = "/usr/local/vamp/haproxy.cfg"
keys = [ "${VAMP_KEY_VALUE_STORE_PATH}" ]
check_cmd = "/usr/local/vamp/validate.sh {{.src}}"
reload_cmd = "/usr/local/vamp/haproxy-reload.sh"
EOT

cp /usr/local/vamp/haproxy.basic.cfg /usr/local/vamp/confd/templates/haproxy.tmpl
cat <<EOT >> /usr/local/vamp/confd/templates/haproxy.tmpl
{{getv "${VAMP_KEY_VALUE_STORE_PATH}"}}
EOT

cat <<EOT > /usr/local/vamp/confd/conf.d/filebeat.toml
[template]
src = "filebeat.tmpl"
dest = "/usr/local/filebeat/filebeat.yml"
EOT

cat <<EOT >> /usr/local/vamp/confd/templates/filebeat.tmpl
filebeat.prospectors:
- input_type: log
  paths:
    - /var/log/haproxy.log
  json.message_key: message
  json.add_error_key: false
  json.keys_under_root: true

output.elasticsearch:
  hosts: ["$VAMP_PERSISTENCE_STORE_CONNECTION"]
  index: "vamp-vga-%{+yyyy.MM.dd}"

path.home: /usr/local/filebeat
path.config: \${path.home}
path.logs: /var/log
EOT

touch /usr/local/vamp/good2go
