# files/etc/service/vamp-gateway-agent/run
#! /bin/bash

handle() { echo "vamp-gateway-agent/run: got signal"; exit; }
trap handle SIGINT

[[ -e /usr/local/vamp/good2go ]] && exit 0

: "${VAMP_GATEWAY_AGENT_LOGO:=TRUE}"

if [ "$VAMP_GATEWAY_AGENT_LOGO" = "TRUE" ] || [ "$VAMP_GATEWAY_AGENT_LOGO" = "1" ]; then
echo "
██╗   ██╗ █████╗ ███╗   ███╗██████╗      ██████╗  █████╗ ████████╗███████╗██╗    ██╗ █████╗ ██╗   ██╗
██║   ██║██╔══██╗████╗ ████║██╔══██╗    ██╔════╝ ██╔══██╗╚══██╔══╝██╔════╝██║    ██║██╔══██╗╚██╗ ██╔╝
██║   ██║███████║██╔████╔██║██████╔╝    ██║  ███╗███████║   ██║   █████╗  ██║ █╗ ██║███████║ ╚████╔╝
╚██╗ ██╔╝██╔══██║██║╚██╔╝██║██╔═══╝     ██║   ██║██╔══██║   ██║   ██╔══╝  ██║███╗██║██╔══██║  ╚██╔╝
 ╚████╔╝ ██║  ██║██║ ╚═╝ ██║██║         ╚██████╔╝██║  ██║   ██║   ███████╗╚███╔███╔╝██║  ██║   ██║
  ╚═══╝  ╚═╝  ╚═╝╚═╝     ╚═╝╚═╝          ╚═════╝ ╚═╝  ╚═╝   ╚═╝   ╚══════╝ ╚══╝╚══╝ ╚═╝  ╚═╝   ╚═╝
"
fi

: "${VAMP_KEY_VALUE_STORE_PATH?not provided.}"
: "${VAMP_KEY_VALUE_STORE_TYPE?not provided.}"
: "${VAMP_KEY_VALUE_STORE_CONNECTION?not provided.}"
: "${VAMP_ELASTICSEARCH_URL?not provided.}"

for try in {1..9} ; do
  echo "vamp-gateway-agent/run: waiting for syslog to come online..."
  [[ -e /var/run/rsyslogd.pid ]] && break || sleep 1
done

if [ ! z "$(ls -A /usr/local/share/ca-certificates)" ]; then
   echo "Adding custom CA certificate"
   update-ca-certificates 2>/dev/null || true
fi

printf "VAMP GATEWAY AGENT: " && cat /usr/local/vamp/version
/usr/bin/confd -version

echo "VAMP_KEY_VALUE_STORE_TYPE       : ${VAMP_KEY_VALUE_STORE_TYPE}"
echo "VAMP_KEY_VALUE_STORE_CONNECTION : ${VAMP_KEY_VALUE_STORE_CONNECTION}"
echo "VAMP_KEY_VALUE_STORE_PATH       : ${VAMP_KEY_VALUE_STORE_PATH}"
echo "VAMP_ELASTICSEARCH_URL          : ${VAMP_ELASTICSEARCH_URL}"

/usr/local/vamp/confd-generate-templates.sh

if [[ "${VAMP_KEY_VALUE_STORE_TYPE}" == "vault" && -n ${VAMP_VAULT_CERTIFICATE_NAME} ]]; then
  curl --cacert ${CONFD_CLIENT_CAKEYS} \
       --cert ${CONFD_CLIENT_CERT} \
       --key ${CONFD_CLIENT_KEY} \
       --request POST --data "{\"name\": \"${VAMP_VAULT_CERTIFICATE_NAME}\"}" -s \
       ${VAMP_KEY_VALUE_STORE_CONNECTION}/v1/auth/cert/login | jq -r .auth.client_token > /usr/local/vamp/token
fi

touch /usr/local/vamp/good2go

if [ "${VAMP_KEY_VALUE_STORE_TYPE}" == "vault" ]; then
  /usr/local/vamp/tokenrenewer.sh &
fi
