# files/etc/service/confd/finish
#! /bin/bash

: "${VAMP_GATEWAY_AGENT_CONFD_INTERVAL:=5}"

sleep $VAMP_GATEWAY_AGENT_CONFD_INTERVAL
exit 0
