# files/etc/service/confd/finish
#! /bin/bash

interval=5

echo "confd/stop: waiting ${interval} second(s) until next poll"
sleep $interval

